# systemd/tasera-scheduling-software.socket
[Unit]
Description=Tasera scheduling software socket

[Socket]
ListenStream=80
NoDelay=true

[Install]
WantedBy=tasera-scheduling-software.service
